# systemd/tapauthd.service
# TapAuth daemon service
[Unit]
Description=TapAuth authentication daemon
Requires=tapauthd.socket
After=dbus.service bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
User=tapauthd
Group=tapauthd
# Socket activation handed via tapauthd.socket
Sockets=tapauthd.socket
# Manage firewall rules as root before/after the service runs
ExecStartPre=+/usr/bin/tapauthd manage-firewall open
ExecStopPost=+/usr/bin/tapauthd manage-firewall close
ExecStart=/usr/bin/tapauthd
Restart=on-failure
StateDirectory=tapauth

# Hardening (adjust as needed if daemon requires more access)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# Allow management of the TOML config file under ProtectSystem=strict
ReadWritePaths=/etc/tapauth
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# Allow the unprivileged service to manage netfilter rules (pairing opens ephemeral TCP ports)
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
