# systemd/polkit-agent-helper@.service.d/tapauth.conf
[Service]
# The '-' prefix suppresses systemd errors if paths are not yet initialized
ReadWritePaths=-/run/tapauthd -/var/log/tapauth
